# supervisord.conf
[program:gifted]
command = node gifted.js
directory = /home/dotcloud/current
autostart=true
autorestart=false
