# supervisord.conf
[program:gifted]
command = node gifted.js
directory = /home/dotcloud/current/server
autostart=true
autorestart=false
