# server/supervisord.conf
[program:ws]
command = node gifted.js
directory = /home/dotcloud/current/
autostart=true
autorestart=true
